# launchd/com.yahoo-auction-watcher.dispatch-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yahoo-auction-watcher.dispatch-check</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/sawadaakira/Projects/MOTHERSHIP/apps/yahoo-auction-watcher/scripts/dispatch-check-if-stale.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/Users/sawadaakira/Projects/MOTHERSHIP/apps/yahoo-auction-watcher</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/sawadaakira/Projects/MOTHERSHIP/apps/yahoo-auction-watcher/logs/dispatch-check.launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/sawadaakira/Projects/MOTHERSHIP/apps/yahoo-auction-watcher/logs/dispatch-check.launchd.err.log</string>
</dict>
</plist>
